# alpine-apache/root/etc/services.d/apache/run
#!/usr/bin/with-contenv sh

#exec /usr/sbin/apachectl -DFOREGROUND;
exec /usr/sbin/httpd -DFOREGROUND;
